# app-emulation/anbox/files/anbox-session-manager.service
[Unit]
Description=Anbox Session Manager

[Service]
ExecStart=/usr/bin/anbox --gles-driver=host session-manager

[Install]
WantedBy=default.target
